# scripts/macos/com.imessagearchiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.imessagearchiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME_DIR__/bin/imessage-archiver</string>
        <string>-config</string>
        <string>__HOME_DIR__/.config/imessage-archiver/config.yaml</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/usr/local/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer> <!-- 4 PM daily -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>__HOME_DIR__</string>
    <key>StandardOutPath</key>
    <string>__HOME_DIR__/Library/Logs/com.imessagearchiver.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/Library/Logs/com.imessagearchiver.err.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
